# supervisord.conf
[include]

files = /etc/supervisor/conf.d/*.conf
[program:app]

command=bash

#directory will be any folder where you wnat supervisor to cd before executing.
directory=/var/www/html

autostart=true

autorestart=true

startretries=3

#user will be anyone you want but make sure that user will have the enough privilage.
user=root

[supervisord]

nodaemon=true



[supervisorctl]
